# conf/etc/profile.d/tourbillon.sh
export PATH=$PATH:/usr/share/python/tourbillon-agent/bin
